# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default, this file does nothing.

# thermal:
echo 111111 > /sys/class/thermal/thermal_zone0/trip_point_0_temp
echo 90000 > /sys/class/thermal/thermal_zone0/trip_point_1_temp
echo 70000 > /sys/class/thermal/thermal_zone0/trip_point_2_temp
echo 60000 > /sys/class/thermal/thermal_zone0/trip_point_3_temp
echo 50000 > /sys/class/thermal/thermal_zone0/trip_point_4_temp

# defaults:
# echo 125000 > /sys/class/thermal/thermal_zone0/trip_point_0_temp
# echo 120000 > /sys/class/thermal/thermal_zone0/trip_point_1_temp
# echo 115000 > /sys/class/thermal/thermal_zone0/trip_point_2_temp
# echo 85000 > /sys/class/thermal/thermal_zone0/trip_point_3_temp
# echo 60000 > /sys/class/thermal/thermal_zone0/trip_point_4_temp

exit 0
