# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default, this file does nothing.

# thermal defaults:
# echo 125000 > /sys/class/thermal/thermal_zone0/trip_point_0_temp
# echo 120000 > /sys/class/thermal/thermal_zone0/trip_point_1_temp
# echo 115000 > /sys/class/thermal/thermal_zone0/trip_point_2_temp
# echo 85000 > /sys/class/thermal/thermal_zone0/trip_point_3_temp
# echo 60000 > /sys/class/thermal/thermal_zone0/trip_point_4_temp

# system reboot threshold:
echo 107000 > /sys/class/thermal/thermal_zone0/trip_point_0_temp

# hot threshold, but pwm1 = 255 for unknown reason:
echo 105000 > /sys/class/thermal/thermal_zone0/trip_point_1_temp

# maximum fan rotation speed (pwm1 = 0) threshold:
echo 70000 > /sys/class/thermal/thermal_zone0/trip_point_2_temp

# minimum fan rotation speed (pwm1 = 96) threshold:
echo 60000 > /sys/class/thermal/thermal_zone0/trip_point_3_temp

# fan stopped (pwm1 = 255) threshold:
echo 50000 > /sys/class/thermal/thermal_zone0/trip_point_4_temp

exit 0
